# p4wnp1.service
[Unit]
Description=P4wnP1-O2 base service (no auto-payload)
After=network-online.target p4wnp1-usb-prep.service
Wants=network-online.target

[Service]
Type=simple
Environment=P4WN_HOME=/opt/p4wnp1
# Minimal: keep process alive, log a banner
ExecStart=/usr/bin/env bash -c 'echo "[P4wnP1] Ready"; sleep infinity'
WorkingDirectory=/opt/p4wnp1
Restart=always

[Install]
WantedBy=multi-user.target
