# systemd/p4wnp1.service
[Unit]
Description=P4wnP1-O2 (USB + OLED + WebUI)
After=network-online.target sys-kernel-config.mount
Wants=network-online.target sys-kernel-config.mount

[Service]
Type=simple
Environment=P4WN_HOME=/opt/p4wnp1
# Prep (no separate unit needed)
ExecStartPre=/sbin/modprobe configfs
ExecStartPre=/bin/mount -t configfs none /sys/kernel/config
ExecStartPre=/sbin/modprobe dwc2
ExecStartPre=/sbin/modprobe libcomposite
# Bring up USB early (auto = RNDIS then ECM if no link)
ExecStartPre=/usr/bin/env python3 /opt/p4wnp1/p4wnctl.py usb auto
ExecStartPre=/usr/bin/env python3 /opt/p4wnp1/p4wnctl.py net share wlan0
# Long-running processes (keep the service alive)
ExecStart=/bin/bash -c "/usr/bin/env python3 /opt/p4wnp1/webui/app.py & echo '[P4wnP1] Ready'; wait -n"
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
